# service/systemd/wifiwatcher.timer
[Unit]
Description=Run wifi watcher every 10 min

[Timer]
OnBootSec=2min
OnUnitActiveSec=10m 

[Install]
WantedBy=timers.target
